# etc/rc.local
#!/bin/sh -e

/usr/local/bin/shared-shared-start-stop start &>/dev/null &
exit 0
